# covalent_ui/os_api/supervisord.conf
[unix_http_server]
file=/tmp/supervisor-osapp.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord-osapp.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-osapp.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-osapp.sock ; use a unix:// URL  for a unix socket


; ***************************
;  APPLICATION
; ***************************

[program:main.py]
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
command=/usr/bin/python3 main.py
stdout_logfile=./logs/main_out.log
stdout_logfile=./logs/main_err.log

autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600

[program:app.py]
logfile_maxbytes=10MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
command=/usr/bin/python3 ../../covalent_ui/app.py
stdout_logfile=./logs/app_out.log
stdout_logfile=./logs/app_err.log

autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
