# supervisord.conf
[group:covalent]
programs=queuer,dispatcher,runner,data,ui_backend

[program:queuer]
directory=/Users/alejandro/Documents/covalent/refactor/queuer
command=python main.py 
autostart=true
autorestart=true
stderr_logfile=/Users/alejandro/Documents/covalent/refactor/queuer/err.log
stdout_logfile=/Users/alejandro/Documents/covalent/refactor/queuer/out.log
numprocs=1
environment=PORT="8000",DATA_OS_SVC_HOST_URI="http://localhost:8003",MQ_CONNECTION_URI="localhost:4222",MQ_DISPATCH_TOPIC="workflow.dispatch"

[program:dispatcher]
directory=/Users/alejandro/Documents/covalent/refactor/dispatcher
command=python main.py 
autostart=true
autorestart=true
stderr_logfile=/Users/alejandro/Documents/covalent/refactor/dispatcher/err.log
stdout_logfile=/Users/alejandro/Documents/covalent/refactor/dispatcher/out.log
numprocs=1
environment=PORT="8001",DATA_OS_SVC_HOST_URI="http://localhost:8003"

[program:runner]
directory=/Users/alejandro/Documents/covalent/refactor/runner
command=python main.py 
autostart=true
autorestart=true
stderr_logfile=/Users/alejandro/Documents/covalent/refactor/runner/err.log
stdout_logfile=/Users/alejandro/Documents/covalent/refactor/runner/out.log
numprocs=1
environment=PORT="8002",DATA_OS_SVC_HOST_URI="http://localhost:8003"

[program:data]
directory=/Users/alejandro/Documents/covalent/refactor/data
command=python main.py 
autostart=true
autorestart=true
stderr_logfile=/Users/alejandro/Documents/covalent/refactor/data/err.log
stdout_logfile=/Users/alejandro/Documents/covalent/refactor/data/out.log
numprocs=1
environment=PORT="8003"

[program:ui_backend]
directory=/Users/alejandro/Documents/covalent/refactor/ui_backend
command=python main.py 
autostart=true
autorestart=true
stderr_logfile=/Users/alejandro/Documents/covalent/refactor/ui_backend/err.log
stdout_logfile=/Users/alejandro/Documents/covalent/refactor/ui_backend/out.log
numprocs=1
environment=PORT="8004",DATA_OS_SVC_HOST_URI="http://localhost:8003"


[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file


[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface


[supervisorctl]
sserverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket