# init/openrc/init.d/covalent
#!/sbin/openrc-run

name="Covalent server"
description=""
command="$COVALENT_ROOT/bin/$COVALENT_PYTHON"
command_args="$COVALENT_ROOT/lib/$COVALENT_PYTHON/site-packages/covalent_ui/app.py --port $COVALENT_SVC_PORT --ignore-migrations"
pidfile=/var/run/covalentd.pid
PATH=$COVALENT_ROOT/bin:$PATH
PYTHONPATH=$COVALENT_ROOT/lib/$COVALENT_PYTHON/site-packages:$PYTHONPATH

depend() {
    need localmount
}

start() {
    ebegin "Starting Covalent Server"
    start-stop-daemon --start --exec \
    env COVALENT_CONFIG_DIR=$COVALENT_CONFIG_DIR \
        COVALENT_EXECUTOR_DIR=$COVALENT_EXECUTOR_DIR \
	COVALENT_RESULTS_DIR=$COVALENT_RESULTS_DIR \
	COVALENT_DATA_DIR=$COVALENT_DATA_DIR \
	COVALENT_DATABASE=$COVALENT_DATABASE \
	COVALENT_CACHE_DIR=$COVALENT_CACHE_DIR \
	COVALENT_LOGDIR=$COVALENT_LOGDIR \
	COVALENT_SVC_PORT=$COVALENT_SVC_PORT \
	COVALENT_SERVER_IFACE_ANY=1 \
    $command \
    --make-pidfile --pidfile $pidfile \
    --background \
    --stderr $COVALENT_LOGDIR/covalentd.log --stdout $COVALENT_LOGDIR/covalentd.log \
    --chdir $COVALENT_ROOT/lib/$COVALENT_PYTHON/site-packages/covalent_ui \
    -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping Covalent Server"
    parent=$(pgrep -F $pidfile)
    child=$(pgrep -P $parent)
    kill -9 $parent $child
    rm -f $pidfile
    eend $?
}
